# systemd/github2forgejo.timer
[Unit]
Description=Synchronize GitHub repos with Forgejo

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
