# systemd/run-media-storage.automount
[Unit]
Description=nfs automount for storage

[Automount]
Where=/run/media/storage
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
