# systemd/run-media-media.automount
[Unit]
Description=nfs automount for media

[Automount]
Where=/run/media/media
TimeoutIdleSec=0

[Install]
WantedBy=multi-user.target
