# systemd/update-botslist.timer
[Unit]
Description=Update list of bots to block in Caddy

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
